# systemd/daily-market-updater.service
[Unit]
Description=Daily Market Update Service (9 AM ET)
After=network.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/workspace
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=utf-8"
Environment="DJANGO_SETTINGS_MODULE=stockscanner_django.settings"

# Run the daily updater once
ExecStart=/usr/bin/python3 /workspace/daily_market_updater.py -once

# Logging
StandardOutput=append:/var/log/daily-market-updater.log
StandardError=append:/var/log/daily-market-updater.error.log

# Timeouts
TimeoutStartSec=1800
TimeoutStopSec=60