# systemd/market-hours-manager.service
[Unit]
Description=Enhanced Market Hours Manager with Integrated Proxy Management
After=network.target postgresql.service mysql.service
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/workspace
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONIOENCODING=utf-8"
Environment="DJANGO_SETTINGS_MODULE=stockscanner_django.settings"
Environment="PREMARKET_START=04:00"
Environment="MARKET_OPEN=09:30"
Environment="MARKET_CLOSE=16:00"
Environment="POSTMARKET_END=20:00"

# Main command
ExecStart=/usr/bin/python3 /workspace/market_hours_manager_enhanced.py

# Restart policy
Restart=always
RestartSec=30

# Logging
StandardOutput=append:/var/log/market-hours-manager.log
StandardError=append:/var/log/market-hours-manager.error.log

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Timeouts
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target